# open-webui/rootfs/etc/services.d/open-webui/run
#!/usr/bin/with-contenv bashio

# Get the configured port
PORT=$(bashio::config 'port')

# Export the port as an environment variable
export PORT="${PORT}"

# Print the configuration
bashio::log.info "Starting Open WebUI with port ${PORT}"

# Process custom environment variables
for var in $(bashio::config 'env_vars|keys'); do
    name=$(bashio::config "env_vars[${var}].name")
    value=$(bashio::config "env_vars[${var}].value")
    bashio::log.info "Setting environment variable ${name}"
    export "${name}=${value}"
done

# Call the original start.sh script
cd /app/backend
exec ./start.sh